# deploy/systemd/cryptotrader-decide.service
[Unit]
Description=CryptoTrader Trading Decisions
After=network.target postgresql.service cryptotrader-sentiment.service

[Service]
Type=oneshot
User=cryptotrader
Group=cryptotrader
WorkingDirectory=/opt/cryptotrader
Environment="PATH=/opt/cryptotrader/.venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/opt/cryptotrader/.env

ExecStart=/opt/cryptotrader/.venv/bin/python /opt/cryptotrader/main.py --task decide
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cryptotrader-decide
