# deploy/systemd/cryptotrader-execute.service
[Unit]
Description=CryptoTrader Order Execution
After=network.target postgresql.service cryptotrader-decide.service

[Service]
Type=oneshot
User=cryptotrader
Group=cryptotrader
WorkingDirectory=/opt/cryptotrader
Environment="PATH=/opt/cryptotrader/.venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/opt/cryptotrader/.env

ExecStart=/opt/cryptotrader/.venv/bin/python /opt/cryptotrader/main.py --task execute
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cryptotrader-execute
